# betting-arb-healthcheck.timer
[Unit]
Description=Run betting-arbitrage health check every 5 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=5min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target