# betting-arb-cleanup.timer
[Unit]
Description=Run betting-arbitrage disk cleanup every 2 hours

[Timer]
OnBootSec=15min
OnUnitActiveSec=2h
AccuracySec=5min
Persistent=true

[Install]
WantedBy=timers.target