# betting-arb-healthcheck.timer
[Unit]
Description=Run betting-arbitrage health check every 10 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=10min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
